# systemd/hscout.timer
[Unit]
Description=Hetzner Hunter
Wants=hscout.service

[Timer]
Unit=hscout.service
OnCalendar=hourly

[Install]
WantedBy=timers.target
